# home/.config/systemd/user/app-syncthingtray@autostart.service.d/override.conf
[Service]
# give tray time to start up
ExecStartPre=sleep 5s
