# udev/rules.d/90-drmwatch.rules
SUBSYSTEM=="drm", RUN+="/bin/bash /home/seckinger/.display.sh"
