# udev/rules.d/90-drmwatch.rules
SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/seckinger/.Xauthority", RUN+="/bin/bash /home/seckinger/.display.sh"
