# systemd/system/certbot-renew.timer
[Unit]
Description=renew Let's Encrypt certificates when due

[Timer]
OnStartupSec=1min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
