# systemd/system/ddns@.service
[Unit]
Description=update ddns

[Service]
Type=simple
ExecStart=/usr/bin/env curl -s http://sync.afraid.org/u/%i/
