# systemd/system/ddns6@.timer
[Unit]
Description=Update ddns v6

[Timer]
Unit=ddns6@%i.service
OnStartupSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
