# systemd/system/ddns@.timer
[Unit]
Description=Update ddns

[Timer]
Unit=ddns@%i.service
OnStartupSec=30sec
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
