# systemd/system/hostblock.timer
[Unit]
Description=update blocked hosts on boot

[Timer]
OnBootSec=10min

[Install]
WantedBy=timers.target
