# systemd/system/hostblock.timer
[Unit]
Description=update blocked hosts 1min after startup

[Timer]
OnStartupSec=1min

[Install]
WantedBy=timers.target
