# systemd/system/hostblock.timer
[Unit]
Description=update blocked hosts 1min after startup

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
