# etc/systemd/system/sshd.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=42337
