# home/.config/autostart/xss-lock.desktop
[Desktop Entry]
Name=XSS Lock
Exec=xss-lock -l -- xss-lock-script
Type=Application
X-GNOME-Autostart-enabled=true
