# home/.config/systemd/user/app-solaar@autostart.service.d/override.conf
[Service]
ExecStart=
ExecStart=bash -c 'sleep 3s && solaar --window=hide'
