# home/.config/systemd/user/app-solaar@autostart.service.d/override.conf
[Service]
ExecStart=
ExecStart=solaar --window=hide
