# home/.config/systemd/user/github-sync.service
[Unit]
Description=fetch all GitHub repos into .github

[Service]
Type=simple
WorkingDirectory=%h/.github
ExecStart=zsh -c "yarn && ./sync.js"
