# home/.config/systemd/user/grive.service
[Unit]
Description=grive sync

[Service]
Type=simple
ExecStart=/usr/bin/grive -p %h/gdrive
