# home/.config/systemd/user/grive.service
[Unit]
Description=grive sync

[Service]
Type=simple
ExecStart=%h/bin/grive -p %h/gdrive
TimeoutSec=30min
