# home/.config/systemd/user/i3-session.path
# https://github.com/i3/i3/issues/5186#issuecomment-2463357316
[Unit]
Description="Waits for i3 to startup before binding a service to it"

[Path]
PathExistsGlob=/run/user/%U/i3/ipc-socket.*

[Install]
WantedBy=graphical-session.target
