# home/.config/systemd/user/npm-global-export.timer
[Unit]
Description=run npm-global-export once a week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
