# home/.config/systemd/user/photoprism-index.timer
[Unit]
Description=Index PhotoPrism originals once per day

[Timer]
OnCalendar=*-*-* 05:15:00
Persistent=false

[Install]
WantedBy=timers.target
