# home/.config/systemd/user/photoprism-index.timer
[Unit]
Description=Index PhotoPrism originals once per day

[Timer]
OnCalendar=*-*-* 04:55:00
Persistent=false

[Install]
WantedBy=timers.target
