# home/.config/systemd/user/pkg-export.service
[Unit]
Description=dump paru package list

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'paru -Qe >%h/docs/tech/pkg-list/%H.txt'
