# home/.config/systemd/user/pkg-export.service
[Unit]
Description=dump yay package list

[Service]
Type=simple
ExecStart=/usr/bin/bash -c 'yay -Qe >%h/docs/tech/pkg-list/%H.txt'
