# home/.config/systemd/user/shell-history-export.service
[Unit]
Description=backup shell history

[Service]
Type=simple
ExecStart=/usr/bin/cp %h/.histfile %h/docs/tech/shell-history/%H.histfile
