# home/.config/systemd/user/shell-history-export.timer
[Unit]
Description=run pkg-export every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
