# home/.config/systemd/user/sway-session.target
[Unit]
Description=sway session is up
Requires=basic.target
BindsTo=graphical-session.target xdg-desktop-autostart.target

RefuseManualStart=yes
RefuseManualStop=yes
StopWhenUnneeded=yes
