# home/.config/systemd/user/vim-plug-export.timer
[Unit]
Description=run vim-plug-export once a week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
