# home/.config/systemd/user/vscode-ext-export.timer
[Unit]
Description=run vscode-ext-export once a week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
